# supervisord.prod.conf
# Production Supervisor Configuration for Tailscale Exit Node
# Enhanced security, logging, and process management

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
username=supervisor
password=change_this_password_in_production

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true  # Run in foreground for Docker
minfds=1024
minprocs=200
user=root
silent=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=supervisor
password=change_this_password_in_production

# Unbound DNS Resolver
[program:unbound]
command=/usr/sbin/unbound -d -c /etc/unbound/unbound.conf
user=root
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stderr_logfile=/var/log/supervisor/unbound.err.log
stdout_logfile=/var/log/supervisor/unbound.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB

# Tailscaled daemon
[program:tailscaled]
command=/usr/local/bin/tailscaled --state=/var/lib/tailscale/tailscaled.state --socket=/var/run/tailscale/tailscaled.sock --port=41641 --tun=kernel
user=root  # Required for network operations
autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=30
stderr_logfile=/var/log/tailscale/tailscaled.err.log
stdout_logfile=/var/log/tailscale/tailscaled.out.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stderr_logfile_backups=5
stdout_logfile_backups=5
redirect_stderr=false
stdout_events_enabled=false
stderr_events_enabled=false

# Tailscale client connection
[program:tailscale-up]
command=/bin/sh -c 'sleep 10 && /usr/local/bin/tailscale up --authkey=${TAILSCALE_AUTHKEY} --hostname=${HOSTNAME:-exit-node} --advertise-exit-node --accept-routes --accept-dns=false --netfilter-mode=on --snat-subnet-routes=false --reset'
user=root
autostart=true
autorestart=false  # Only run once
startretries=3
startsecs=5
stopwaitsecs=10
stderr_logfile=/var/log/tailscale/tailscale-up.err.log
stdout_logfile=/var/log/tailscale/tailscale-up.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3

# Nginx web server
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root  # Required to bind to privileged ports
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stderr_logfile=/var/log/nginx/nginx.err.log
stdout_logfile=/var/log/nginx/nginx.out.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stderr_logfile_backups=5
stdout_logfile_backups=5

# Fail2ban for intrusion prevention
[program:fail2ban]
command=/usr/bin/fail2ban-server -f
user=root
autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=30
stderr_logfile=/var/log/supervisor/fail2ban.err.log
stdout_logfile=/var/log/supervisor/fail2ban.out.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stderr_logfile_backups=3
stdout_logfile_backups=3

# Log rotation process
[program:logrotate]
command=/usr/sbin/logrotate -s /var/lib/logrotate.status /etc/logrotate.conf
user=root
autostart=false
autorestart=false
startsecs=0

# Self-healing: iptables manager
[program:iptables-manager]
command=/app/iptables-manager.sh monitor
user=root
autostart=true
autorestart=true
startretries=5
startsecs=15
stopwaitsecs=30
stderr_logfile=/var/log/tailscale/iptables-manager.err.log
stdout_logfile=/var/log/tailscale/iptables-manager.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3

# Self-healing: Connection monitor
[program:connection-monitor]
command=/app/connection-monitor.sh monitor
user=root
autostart=true
autorestart=true
startretries=5
startsecs=25
stopwaitsecs=30
stderr_logfile=/var/log/tailscale/connection-monitor.err.log
stdout_logfile=/var/log/tailscale/connection-monitor.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
